# time-leak.service
[Unit]
Description=Time Leak Backend Service
Documentation=https://github.com/your-repo/time-leak
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/time-leak-back-end

# Environment variables
Environment="APP_ADDR=:8080"
Environment="DB_PATH=/home/time-leak-back-end/data"
Environment="DB_NAME=timeleak.db"
Environment="PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Run using go run
ExecStart=/usr/local/go/bin/go run ./cmd

# Restart policy
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (adjust as needed)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
